# scripts/com.ai-daily.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-daily.daily-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/wangjd/clawd/ai-daily/scripts/daily_sync.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL</key>
        <string>postgresql://postgres:postgres@localhost:5432/ai_daily</string>
        <key>TAVILY_API_KEY</key>
        <string>tvly-your-tavily-key</string>
        <key>MINIMAX_API_KEY</key>
        <string>your-minimax-key</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/wangjd/clawd/ai-daily/logs/daily_sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wangjd/clawd/ai-daily/logs/daily_sync_error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>